# templates/default/mailcatcher.upstart.conf
#!upstart
description "Mailcatcher SMTP dev tool"

start on startup
stop on shutdown

respawn

exec mailcatcher