# templates/default/mailcatcher.upstart.conf.erb
#!upstart
description "Mailcatcher SMTP dev tool"

start on startup
stop on shutdown

exec mailcatcher --http-ip="<%= node['mailcatcher']['http_ip'] %>"
